# runit/remote-api-server/run
#!/bin/bash
export FLASK_RUN_PORT=8000
export FLASK_APP=server
exec flask run --host=0.0.0.0
